# main/lxdm/lxdm.service
[Unit]
Description=LXDE Display Manager

[Service]
ExecStart=/usr/sbin/lxdm

[Install]
WantedBy=graphical.target
